# fs/nmbd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM 
"/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='Vivitec:nmbd'>

<service
         name='application/nmbd'
         type='service'
         version='1'>

          <create_default_instance enabled='true' />

         <single_instance/>

         <dependency
                 name='name-services'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/milestone/name-services' />
         </dependency>

         <dependency
                 name='local-filesystems'
                 type='service'
                 grouping='require_all'
                 restart_on='none'>
                 <service_fmri value='svc:/system/filesystem/local' />
         </dependency>

  <dependency name='config_data'
    grouping='require_all'
    restart_on='restart'
    type='path'>
    <service_fmri
        value='file://localhost/opt/local/etc/samba/smb.conf' />
  </dependency>


         <exec_method
                 type='method'
                 name='start'
                 exec='/usr/local/svc/method/nmbd %m'
                 timeout_seconds='60' />

         <exec_method
                 type='method'
                 name='stop'
                 exec='/usr/local/svc/method/nmbd %m'
                 timeout_seconds='60' />

         <stability value='Unstable' />

         <template>
                 <common_name>
                         <loctext xml:lang='C'>Samba
                         </loctext>
                 </common_name>
         </template>
</service>

</service_bundle>

